# plc-gateway.service
[Unit]
Description=PLC Edge Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/plc-diagnostics
EnvironmentFile=/home/pi/plc-diagnostics/.env
ExecStart=/home/pi/plc-diagnostics/venv/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 5000
Restart=always
RestartSec=10
TimeoutStopSec=20
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plc-gateway

[Install]
WantedBy=multi-user.target
